# conf/supervisor.conf
[program:vae]
environment=FLASK_APP="/home/dmurphy/vaedbl/vaedbl.py"
command=flask run --host=0.0.0.0 --port=8080
