# conf/supervisor.conf
[program:vae]
environment=FLASK_APP="<path>/vae_flask/flask.py"
command=flask run --host=0.0.0.0
