# conf/supervisor.conf
[program:vae]
environment=FLASK_APP="<path>/vaedbl/vaedbl.py"
command=flask run --host=0.0.0.0
